# etc/supervisord.d/ssh.ini
[program:sshd]
command=/usr/sbin/sshd -D
stopwaitsecs=60
prioirty=10
autorestart=true
