# etc/supervisord.d/tcollector.ini
[program:tcollector]
command=/opt/tcollector/bin/tcollector.sh
startsecs=60
stopwaitsecs=60
prioirty=10
autorestart=true
