# etc/supervisord.d/tsddrain.ini
[program:tsd_drain]
command=/opt/tcollector/bin/tsddrain.sh
stopwaitsecs=60
prioirty=10
autorestart=true
autostart=false
